# config/docker/docker.service.d/http-proxy.conf
[Service]
Environment="HTTP_PROXY=http://localhost:8800"
Environment="HTTPS_PROXY=http://localhost:8800"
Environment="NO_PROXY=localhost,127.0.0.1,.local"
