# config/docker/docker.service.d/http-proxy.conf
[Service]
EnvironmentFile=/etc/envs/http_proxy.env
