# config/systemd/system/xray.service.d/local.conf
[Service]
User=root
ExecStart=
ExecStart=/usr/local/bin/xray -confdir /etc/xray/
Environment=
Environment=XRAY_LOCATION_ASSET=/usr/share/v2ray-rules-dat
