# config/xray-local.service
[Service]
ExecStart=
ExecStart=/usr/local/bin/xray -confdir /usr/local/etc/xray/
